# gnomeos/yocto/recipies-core/initscripts/initscripts-1.0/hostname.sh
### BEGIN INIT INFO
# Provides:          hostname
# Required-Start:
# Required-Stop:
# Default-Start:     S
# Default-Stop:
# Short-Description: Set hostname based on /etc/hostname
### END INIT INFO

if test -f /etc/hostname
then
	hostname -F /etc/hostname
fi
